# packages/rootfs-dev/console.conf
description "getty console"

start on sys-init

respawn
exec /sbin/getty -n -l /bin/sh console 115200 vt100
