# packages/rootfs-dev/network-init.conf
description "Network setup on init"

start on sys-init

task

exec /etc/scripts/network-init.sh