# packages/rootfs-dev/startup.conf
description "System boot up"

start on startup

task

emits sys-init

exec /etc/scripts/startup.sh
