# packages/rootfs-dev/tty2.conf
description "getty"

start on sys-init

respawn
exec /sbin/getty -n -l /bin/sh tty2 115200 vt100
