# packages/rootfs-vapor/init/console.conf
description "console"

start on sys-boot

console owner

respawn

script
    if [ -x /etc/scripts/console.sh ]; then
        exec /etc/scripts/console.sh
    else
        exec /bin/bash -l
    fi
end script
