# packages/rootfs-vapor/init/klogd.conf
description "Kernel logging service"

start on sys-boot

respawn

exec /sbin/klogd -n