# packages/rootfs-vapor/init/mdev.conf
description "Scan hardware"

start on mdev-scan

task

exec /sbin/mdev -s
