# packages/rootfs-vapor/init/ovs-dbinit.conf
description "Open vSwitch database initialization"

start on started ovs-dbserver

task

exec ovs-vsctl --no-wait init
