# packages/rootfs-vapor/init/ovs-vswitchd.conf
description "Open vSwitch daemon"

start on started ovs-dbserver

respawn

exec ovs-vswitchd --log-file=/var/log/openvswitch/ovs-vswitchd.log unix:/var/run/openvswitch/db.sock
