# packages/rootfs-vapor/init/startup.conf
description "System boot up"

start on startup

task

script
    /etc/scripts/startup.sh
    initctl emit sys-boot
    initctl emit sys-init
end script
