# packages/rootfs-vapor/init/tty1.conf
description "getty tty1"

start on sys-init

respawn
exec /sbin/getty -n -l /etc/scripts/terminal.sh tty1 115200
