# packages/rootfs-boot/init/tty2.conf
description "getty tty2"

start on sys-init

respawn
exec /sbin/getty -n -l /etc/scripts/terminal.sh tty2 115200 vt100
