# packages/rootfs-boot/init/tty3.conf
description "getty tty3"

start on sys-init

respawn
exec /sbin/getty -n -l /etc/scripts/terminal.sh tty3 115200
