# packages/rootfs-boot/net-if-changed.conf
description "Network Interface state changed"

start on net-if-changed

task

exec /etc/scripts/net-if-changed.sh $NETIF