# packages/rootfs-boot/startup.conf
description "System boot up"

start on startup

task

exec /etc/scripts/startup.sh
